# tests/nested/dir/files/bar.service
[Service]
ExecStart=/bin/sleep 3600
ExecReload=/bin/touch /tmp/bar-reloaded

[Install]
WantedBy=multi-user.target
